# Automation.service
[Unit]
Description=Firefox Web Browser
After=network.target

[Service]
ExecStart=/usr/bin/firefox
Restart=always
RestartSec=3
User=your_username

[Install]
WantedBy=multi-user.target
